# ollama.service
[Unit]
Description=Ollama Service
After=network.target

[Service]
ExecStart=/usr/local/bin/ollama serve --host $OLLAMA_HOST
Environment=OLLAMA_HOST=localhost
Restart=always
User=ubuntu

[Install]
WantedBy=multi-user.target
